# zhengtu.service
[Unit]
Description=Run zhengtu tool

[Service]
ExecStart=%h/.zhengtu/bin/zhengtu

